# remote/supervisord.xrdp.conf
# SPDX-FileCopyrightText: Copyright DB InfraGO AG and contributors
# SPDX-License-Identifier: Apache-2.0

[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
user=techuser
autorestart=true
stdout_logfile=/var/log/session/%(program_name)s.stdout.log
stderr_logfile=/var/log/session/%(program_name)s.stderr.log

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
user=techuser
autorestart=true
environment=DISPLAY=":10"
stdout_logfile=/var/log/session/%(program_name)s.stdout.log
stderr_logfile=/var/log/session/%(program_name)s.stderr.log
